# init/pkgstats.timer
[Unit]
Description=Weekly pkgstats submission

[Timer]
OnCalendar=weekly
RandomizedOffsetSec=1 week
Persistent=true
