# init/pkgstats.timer
[Unit]
Description=Weekly pkgstats submission

[Timer]
OnCalendar=weekly
AccuracySec=1d
RandomizedDelaySec=7d
Persistent=true
